# files/laptop/etc/udev/rules.d/thunderbolt3-dock.rules
ACTION=="add", SUBSYSTEM=="thunderbolt", RUN+="/opt/ytlaces/dock.sh dock"
ACTION=="remove", SUBSYSTEM=="thunderbolt", RUN+="/opt/ytlaces/dock.sh undock"
