# files/desktop/mnt-data.mount
[Mount]
What=4EAC6E89010CB206
Where=/mnt/data
